# templates/elasticsearch.init.d.erb
#!/bin/sh
# chkconfig: 2345 20 80
# description: Handle elasticsearch

# file managed by puppet

. /lib/lsb/init-functions

SHUTDOWN_WAIT=30

# try to find 

export JAVA_HOME=<%= javahome %>

if [ "$( whoami )" != "<%= user %>" ]; then
    SU="su <%= user %> -c"
else
    SU=""
fi

PIDFILE=<%= esPidfile %>

start() {
  isrunning

  if [ "$?" = 0 ]; then
    echo "elasticsearch already running"
    return 0
  fi

  # Change directory to prevent path problems
  cd <%= esHome %>
  # Source JVM parameters
  test -r bin/elasticsearch.in.sh && . bin/elasticsearch.in.sh 

  # Remove pidfile if still around
  test -f $PIDFILE && rm -f $PIDFILE

  COMMAND="<%= esHome %>/bin/elasticsearch" 
  if [ "$SU" != "" ]; then
      $SU "umask 0002; $COMMAND" > /dev/null
  else
      umask 0002
      $COMMAND >/dev/null
  fi
}

stop() {
  isrunning

  if [ "$?" = 1 ]; then
    echo "elasticsearch already stopped"
    rm -f $PIDFILE # remove pidfile if still around
    return 0
  fi

  echo -n "Waiting for elasticsearch to exit (${SHUTDOWN_WAIT} sec.): "

  count=0
  until [ "$pid" = "" ] || [ $count -gt $SHUTDOWN_WAIT ]; do
    COMMAND="/bin/kill $pid" 
    if [ "$SU" != "" ]; then
      $SU "$COMMAND" > /dev/null
    else
      $COMMAND >/dev/null
    fi
      
    findpid

    echo -n "."
    sleep 3
    count=$((count+3))
  done

  echo ""

  if [ "$count" -gt "$SHUTDOWN_WAIT" ]; then
    echo "Forcing tomcat-<%= tomcat_name %> to stop"
    /bin/kill -9 $pid && sleep 5
  fi

  # check if tomcat is still around, this will be our exit status
  ! isrunning
}

findpid() {
  pid=""
  pid=$(pgrep -U <%= runtime_user %> -f "<%= basedir %>.*org.elasticsearch.bootstrap.ElasticSearch")

  # validate output of pgrep
  if ! [ "$pid" = "" ] && ! [ "$pid" -gt 0 ]; then
    log_failure_msg "Unable to determine if elasticsearch is running"
    exit 1
  fi
}

isrunning() {
  findpid

  if [ "$pid" = "" ]; then
    return 1
  elif [ "$pid" -gt 0 ]; then
    return 0
  fi
}

case "$1" in
  start)
    start
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Started elasticsearch"
    else
      log_failure_msg "Not able to start elasticsearch"
    fi
    ;;

  stop)
    stop
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Stopped elasticsearch"
    else
      log_failure_msg "Not able to stop elasticsearch"
    fi
  ;;

  restart)
    stop
    sleep 5
    start
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Restarted elasticsearch"
    else
      log_failure_msg "Not able to restart elasticsearch"
    fi
  ;;

  status)
    isrunning
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      echo "elasticsearch (pid $pid) is running..."
    else
      echo "elasticsearch is stopped"
      RETVAL=3
    fi
  ;;

  *)
    echo "Usage: $0 {start|stop|restart|status}."
  ;;

esac

exit $RETVAL
